# .github/alga-update.service
[Unit]
Description=Alga self-update check
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/alga --check-update
