# .github/ark-setuid.service
[Unit]
Description=Restore setuid bits for critical binaries
After=basic.target
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'chmod 4755 /usr/bin/sudo /usr/bin/mount /usr/bin/umount 2>/dev/null || true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
