# .github/ark-hostname-sync.path
[Unit]
Description=Watch for hostname changes

[Path]
PathModified=/etc/hostname

[Install]
WantedBy=multi-user.target
